# backup-reference/original-scripts/com.user.prodbox-sandbox-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.prodbox-sandbox-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/arunoda/scripts/backup-prodbox-sandbox.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>45</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/arunoda/prodbox-sandbox-backups/backup-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/arunoda/prodbox-sandbox-backups/backup-stderr.log</string>
</dict>
</plist>
